# ops/db-ha/standby/pg-primary-tunnel.service
[Unit]
Description=Maintain standby SSH tunnel to primary Postgres
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/pg-primary-tunnel
ExecStart=/usr/local/bin/maintain-standby-tunnel.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
